# harbour-dyncal/usr/share/harbour-dyncal/harbour-dyncal.timer
[Unit]
Description=Timer for DynCal.

[Timer]
OnBootSec=0
OnCalendar=00:00
OnUnitActiveSec=12h
Persistent=True

[Install]
WantedBy=timers.target
