# usr/share/harbour-dyncal/harbour-dyncal.service
[Unit]
<<<<<<< HEAD
Description=Change calendar icon.
=======
Description=Run the script and change calendar icon.
>>>>>>> c6234cee16cd4c6023e5d96a95d9f3e6b04848a7

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=/usr/share/harbour-dyncal/harbour-dyncal.sh

[Install]
<<<<<<< HEAD
WantedBy=multi-user.target
=======
WantedBy=multi-user.target
>>>>>>> c6234cee16cd4c6023e5d96a95d9f3e6b04848a7
